# systemd/hostapd.service.d/ip-assignment.conf
[Service]
# Assign IP to wlan0 before hostapd starts
ExecStartPre=/opt/pixelpi/scripts/assign-wlan0-ip.sh
